# 1-20March-address-service-optimized-WorkingGudJsonBodyBkp/config/address-service.service
[Unit]
Description=Address Normalization Service
Documentation=https://github.com/your-org/address-service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=addrservice
Group=addrservice
WorkingDirectory=/opt/address-service

ExecStart=/opt/address-service/bin/address-service \
    --config /etc/address-service/config.json

# Graceful shutdown: wait for in-flight batch requests
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

# Auto-restart on failure
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# File descriptor limits (needed for high-concurrency)
LimitNOFILE=65535
LimitNPROC=65535

# Environment
Environment=LIBPOSTAL_DATA_DIR=/opt/libpostal/data
Environment=CACHE_MAX_SIZE=2000000
Environment=LLM_ENABLED=false
Environment=LOG_LEVEL=info

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/address-service /opt/address-service/logs
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=address-service

[Install]
WantedBy=multi-user.target
